# profile.d/acq2106_hts.sh
export HAPI=/home/dt100/PROJECTS/acq400_hapi
export AFHBA=/home/dt100/PROJECTS/AFHBA404
export HTSDATA=/mnt/datastream/ACQ400DATA
# fast
export HTSARCHIVE=/mnt/datastream/SHOTS
# slow
#export HTSARCHIVE=/home/dt100/DATA/MIR


